# lxdm/.footprint
drwxr-xr-x	mitch/users	etc/
drwxr-xr-x	mitch/users	etc/lxdm/
-rwxr-xr-x	mitch/users	etc/lxdm/LoginReady
-rwxr-xr-x	mitch/users	etc/lxdm/PostLogin
-rwxr-xr-x	mitch/users	etc/lxdm/PostLogout
-rwxr-xr-x	mitch/users	etc/lxdm/PreLogin
-rwxr-xr-x	mitch/users	etc/lxdm/PreReboot
-rwxr-xr-x	mitch/users	etc/lxdm/PreShutdown
-rwxr-xr-x	mitch/users	etc/lxdm/Xsession
-rw-r-----	mitch/users	etc/lxdm/lxdm.conf
drwxr-xr-x	mitch/users	etc/pam.d/
-rw-r--r--	mitch/users	etc/pam.d/lxdm
drwxr-xr-x	mitch/users	etc/rc.d/
-rwxr-xr-x	mitch/users	etc/rc.d/lxdm
drwxr-xr-x	mitch/users	usr/
drwxr-xr-x	mitch/users	usr/bin/
-rwxr-xr-x	mitch/users	usr/bin/lxdm-config
drwxr-xr-x	mitch/users	usr/libexec/
-rwxr-xr-x	mitch/users	usr/libexec/lxdm-greeter-gdk
-rwxr-xr-x	mitch/users	usr/libexec/lxdm-greeter-gtk
-rwxr-xr-x	mitch/users	usr/libexec/lxdm-numlock
-rwxr-xr-x	mitch/users	usr/libexec/lxdm-session
drwxr-xr-x	mitch/users	usr/sbin/
-rwxr-xr-x	mitch/users	usr/sbin/lxdm
-rwxr-xr-x	mitch/users	usr/sbin/lxdm-binary
drwxr-xr-x	mitch/users	usr/share/
drwxr-xr-x	mitch/users	usr/share/lxdm/
-rw-r--r--	mitch/users	usr/share/lxdm/config.ui
-rw-r--r--	mitch/users	usr/share/lxdm/lxdm.glade
